# templates/mongoshard.conf.erb
# Ubuntu upstart file at /etc/init/mongoshard.conf
# Managed by Puppet, DO NOT TOUCH!


start on runlevel [2345]
stop on runlevel [06]

script
  ENABLE_MONGOSHARD="yes"
  if [ -f /etc/default/mongoshard ]; then . /etc/default/mongoshard; fi
  if [ "x$ENABLE_MONGOSHARD" = "xyes" ]; then exec start-stop-daemon --start --quiet --chuid mongodb --pidfile /var/run/mongoshard-<%= name %> --exec  /usr/bin/mongod -- --shardsvr --dbpath <%= prefix %>/<%= name %> --port <%= port %> --replSet <%= name %> --logpath /var/log/mongodb/mongoshard-<%= name %>.log --logappend --cpu; fi
end script
